# systemd/user/gpg-refresh.timer
[Unit]
Description=run gpg-refresh on boot

[Timer]
OnBootSec=10min

[Install]
WantedBy=timers.target
