# systemd/user/keepass.service
[Unit]
Description=Start KeePass password manager

[Service]
ExecStart=/bin/keepass

[Install]
WantedBy=wm.target
