# systemd/user/keepass.service
[Unit]
Description=Start KeePass password manager

[Service]
Type=simple
ExecStart=/bin/keepass
Restart=always

[Install]
WantedBy=wm.target
