# systemd/user/wm.target
[Unit]
Description=WM is loaded and available
After=graphical.target
