# systemd/user/keepass.service.d/wants.conf
[Unit]
Wants=dropbox.service ssh-agent.service
After=dropbox.service ssh-agent.service


# vim: filetype=systemd
